# getTime.cron
*/5 * * * * python /home/jakob/Projects/Googlemaps/getTime.py
