# getTime.cron
*/5 * * * * python getTime.py >> results.txt
